# machine-spirit.service
[Unit]
Description=Machine Spirit AI Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/aaron/self-learning-ai/evolve_ai.py
WorkingDirectory=/home/aaron/self-learning-ai
Restart=always
User=aaron
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
