# systemd/machinespirit-reflect.timer
[Unit]
Description=Run MachineSpirit reflection daily

[Timer]
OnCalendar=*-*-* 03:05:00
Persistent=true

[Install]
WantedBy=timers.target
