# systemd/machinespirit-ui.service.in
[Unit]
Description=MachineSpirit UI (FastAPI/Uvicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=@REPO_DIR@
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=%h/.config/machinespirit/secrets.env
ExecStart=@REPO_DIR@/.venv/bin/python -m uvicorn ms_ui:app --host 0.0.0.0 --port 8020
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
