# systemd/machinespirit-autoapply.service
[Unit]
Description=MachineSpirit Auto-Apply (apply latest proposal + selftest + auto-commit)
After=network-online.target machinespirit-api.service machinespirit-ui.service
Wants=network-online.target machinespirit-api.service machinespirit-ui.service

[Service]
Type=oneshot
WorkingDirectory=%h/self-learning-ai
Environment=PYTHONUNBUFFERED=1
Environment=MS_AUTOPUSH=0
SyslogIdentifier=machinespirit-autoapply
StandardOutput=journal
StandardError=journal

# upgrades can take time; avoid premature timeouts
TimeoutStartSec=3600

ExecStart=/bin/bash -lc 'cd %h/self-learning-ai && ./scripts/autoapply_runner.sh'
