# systemd/machinespirit-autopropose.service
[Unit]
Description=MachineSpirit auto-propose (creates pending proposals only)

[Service]
Type=oneshot
WorkingDirectory=%h/self-learning-ai
ExecStart=/usr/bin/python3 %h/self-learning-ai/scripts/auto_propose.py
