# systemd/machinespirit-webqueue.service
[Unit]
Description=MachineSpirit webqueue learner (headless)

[Service]
Type=oneshot
WorkingDirectory=%h/self-learning-ai
ExecStart=/usr/bin/python3 %h/self-learning-ai/brain.py --webqueue --limit 3
